# bin/go1_script.conf
[program:go1_script]
command=/var/www/bin/go1.sh
autostart=true
autorestart=true
stderr_logfile=/var/www/logs/go1.err.log
stdout_logfile=/var/www/logs/go1.out.log